# src/roles/install_local_backup_copy_service/templates/copy_backup.service
[Unit]
Description=Copies newest file to destination

[Service]
ExecStart=/usr/local/bin/copy_backup.sh

[Install]
WantedBy=multi-user.target